# roles/kodi_backup_service/templates/kodi_backup.service.j2
[Unit]
Description = Kodi Automated Backups - Service

[Service]
User = pi
Group = pi
Type = oneshot
ExecStart = {{ kodi_dot_path }}/kodi_backup.sh
